# templates/systemd.service.erb
[Service]
ExecStart=
ExecStart=/usr/sbin/varnishd -a <%= @listen %>:<%= @listen_port %> -T <%= @admin_listen %>:<%= @admin_port %> -f <%= @vcl_conf %> -S <%= @secret_file %> -s <%= @storage_type %>,<%= @storage_size %>
